# systemd/alfred-seo-cwv-sync.service
# systemd/alfred-seo-cwv-sync.service
[Unit]
Description=Alfred SEO — daily CWV sync
After=network-online.target

[Service]
Type=oneshot
User=aialfred
WorkingDirectory=/home/aialfred/alfred
Environment=PYTHONPATH=/home/aialfred/alfred
ExecStart=/home/aialfred/alfred/venv/bin/python /home/aialfred/alfred/scripts/seo_cwv_sync.py
EnvironmentFile=-/home/aialfred/alfred/config/.env
StandardOutput=journal
StandardError=journal
